# ansible/roles/check/files/checkupdate.service
[Unit]
Description=Check updater

[Service]
Type=oneshot
ExecStart=/usr/bin/bash /srv/check.torproject.org/cpexits.sh
