# ansible/roles/check/files/checkupdate.timer
[Unit]
Description=Run check updater every 30 mins

[Timer]
OnBootSec=5min
OnUnitActiveSec=30min

[Install]
WantedBy=timers.target
